# crontab
# nomor 1
14 14 14 2 5 ~/modul_1/soal1.sh
# nomor 4
@hourly ~/modul_1/soal4.sh
# nomor 5
2-30/6 * * * * ~/modul_1/soal5.sh
